# systemd/auctions-updater.service
[Unit]
Description=Periodically pull auction data 

[Service]
Type=simple
User=amas
ExecStart=

[Install]
WantedBy=multi-user.target
